# bat-notification.service
[Unit]
Description=Battery notification

[Service]
ExecStart=/usr/local/bin/battery-notification.py

[Install]
WantedBy=multi-user.target
